# readersv3/apps/barcodeprinter/build/linux/barcodeprinter.service
[Unit]
Description=WiseMED readersv3 Barcode Printer background service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/barcodeprinter
ExecStart=/opt/barcodeprinter/BarcodePrinter -config deployments/config.yaml -headless
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
