# readersv3/apps/generic-test-reader/build/linux/generic-test-reader.service
[Unit]
Description=WiseMED readersv3 Generic Test Reader background service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/generic-test-reader
ExecStart=/opt/generic-test-reader/generic-test-reader -config deployments/config.yaml -headless
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
