# readersv3/apps/shimatzu-generic-reader/build/linux/shimatzu-generic-reader.service
[Unit]
Description=WiseMED readersv3 SHIMATZU Generic Reader background service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/shimatzu-generic-reader
ExecStart=/opt/shimatzu-generic-reader/shimatzu-generic-reader -config deployments/config.yaml -headless
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
