# readersv3/apps/labnovation-ld-560/build/linux/labnovation-ld-560.service
[Unit]
Description=WiseMED readersv3 Labnovation Ld 560 background service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/labnovation-ld-560
ExecStart=/opt/labnovation-ld-560/labnovation-ld-560 -config deployments/config.yaml -headless
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
