# readersv3/apps/cary60-uvvis-reader/build/linux/cary60-uvvis-reader.service
[Unit]
Description=WiseMED readersv3 Cary60 UVVIS Reader background service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cary60-uvvis-reader
ExecStart=/opt/cary60-uvvis-reader/Cary60-UVVIS-Reader -config deployments/config.yaml -headless
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
