# readersv3/apps/update-server/build/macos/eu.wisemed.readersv3.update.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>eu.wisemed.readersv3.update.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/update-server/Update_Server</string>
    <string>-config</string>
    <string>deployments/config.yaml</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/usr/local/update-server</string>
  <key>StandardOutPath</key>
  <string>/var/log/update-server.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/update-server.log</string>
</dict>
</plist>
