# readersv3/apps/update-server/build/linux/update-server.service
[Unit]
Description=WiseMED readersv3 Update Server background service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/update-server
ExecStart=/opt/update-server/Update_Server -config deployments/config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
